# vpn-split/.config/systemd/user/novpn-anchor.service
[Unit]
Description=VPN bypass cgroup anchor

[Service]
Slice=novpn.slice
ExecStart=/usr/bin/sleep infinity
Restart=always

[Install]
WantedBy=default.target
