# vpn-split/.config/systemd/user/novpn.slice
[Unit]
Description=VPN bypass cgroup

[Slice]
